# packaging/postinst
#!/bin/sh
set -e

. /usr/share/debconf/confmodule

case "$1" in
    configure)
        # Get the org ID from debconf
        db_get hardpointd/org_id
        ORG_ID="$RET"

        # Create config directory
        mkdir -p /etc/hardpointd

        # Create key directory
        mkdir -p /var/lib/hardpointd

        # Write config file
        cat > /etc/hardpointd/config.yaml <<EOF
org_id: ${ORG_ID}
EOF

        # Set proper permissions
        chmod 600 /etc/hardpointd/config.yaml
        chmod 755 /var/lib/hardpointd

        # Create a systemd service user
        if ! getent passwd hardpoint >/dev/null; then
            useradd --system --home-dir /var/lib/hardpointd --shell /bin/false hardpoint
        fi

        chown hardpoint:hardpoint /etc/hardpointd/config.yaml
        chown hardpoint:hardpoint /var/lib/hardpointd

        # Reload systemd daemon
        systemctl daemon-reload
        
        # Enable and start service
        systemctl enable hardpointd.service
        systemctl start hardpointd.service || true
        ;;
esac

# Clean up debconf database
db_stop

exit 0